# piwegro_reload.timer
[Unit]
Description=Piwegro reloading service timer

[Timer]
Unit=piwegro_reload.service
OnCalendar=*-*-* 4:00:00

[Install]
WantedBy=timers.target